# .local/projects/printfolder/etc-systemd-system/printchanged.service
[Unit]
Description=Execute print job on pdfs

[Service]
Type=simple
Restart=no
ExecStart=/print/printall.sh

[Install]
WantedBy=multi-user.target