# .local/projects/printfolder/etc-systemd-system/printchanged.path
[Unit]
Description=Execute print job on pdfs

[Path]
PathChanged=/print/monitor
Unit=printchanged.service

[Install]
WantedBy=multi-user.target